# supervisord/api-prdv1.ini
[program:api]
environment=LANG="zh_CN.UTF-8",LANGUAGE="zh_CN:zh:en_US:en",LC_ALL=zh_CN.UTF-8
command=gunicorn -w 4 -t 3000 -b 0.0.0.0:8000 run:app  # /usr/local/bin/gunicorn --worker-class eventlet -w 2 -t 10000 -b 0.0.0.0:8000 run:app
directory=/data/code
stopwaitsecs=0
autostart=true
autorestart=true
stdout_logfile=/data/log/api.log
stderr_logfile=/data/log/api_err.log